# dxx/crontab
0 13 * * 1 /usr/bin/python3 /home/finsh.py > /home/dxx/DxxLog/finsh.log 2>&1 &
